# tasks/mastodon-backup.timer
[Unit]
Description=Mastodon backup timer

[Timer]
OnCalendar=Sun *-*-* 00:00:01
Persistent=true

[Install]
WantedBy=timers.target
