# tasks/mastodon-backup.timer
[Unit]
Description=Mastodon backup timer

[Timer]
OnCalendar=0 1 * * *
Persistent=true

[Install]
WantedBy=timers.target
