# tasks/mastodon-backup.timer
[Unit]
Description=Mastodon backup timer

[Timer]
OnCalendar=*-*-* 4:00:00
Persistent=true

[Install]
WantedBy=timers.target
