# tasks/mastodon-cleanup.timer
[Unit]
Description=Mastodon cleanup timer

[Timer]
OnCalendar=Sat *-*-* 00:00:01
Persistent=true

[Install]
WantedBy=timers.target
