# assets/mastodon-es-healthcheck.timer
[Unit]
Description=ElasticSearch healthcheck timer

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
