# assets/mastodon-backup.service
[Unit]
Description=Mastodon backup

[Service]
Type=oneshot
User=mastodon
Group=mastodon
ExecStart=/home/mastodon/.local/bin/borgmatic -v 2 create
